# vfrb.service.in
[Unit]
Description=VirtualFlightRadar-Backend

[Service]
Type=simple
ExecStart=%VFRB_EXEC_PATH% -c %VFRB_CONF_PATH%

[Install]
WantedBy=multi-user.target
