# service/vfrb.service.d/vfrb.conf
[Service]
Environment="LD_LIBRARY_PATH=%BOOST_LIBS_PATH%:$LD_LIBRARY_PATH"
